# rules/windows/process_creation/susp_schtasks_persistence.yml
title: Scheduled Task Persistence
id: b2c3d4e5-f6a7-8901-bcde-234567890123
status: experimental
description: Detects creation of scheduled tasks from uncommon parent processes or with suspicious names, indicative of persistence mechanisms
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks
author: Rudraksh Dixit
date: 2026-06-21
tags:
    - attack.persistence
    - attack.t1053.005
    - attack.privilege_escalation
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains|all:
            - '/create'
            - '/tn'
    filter_system:
        ParentImage|endswith:
            - '\svchost.exe'
            - '\services.exe'
    filter_admin:
        CommandLine|contains:
            - 'OneDrive'
            - 'GoogleUpdate'
            - 'Adobe'
            - 'Microsoft\Windows'
    condition: selection and not filter_system and not filter_admin
falsepositives:
    - Legitimate software installation creating scheduled tasks
    - Windows built-in scheduled task creation
    - Admin creating maintenance tasks
level: medium
